# templates/etc/init/weave.conf.erb
description "Weave Docker Network Service"
start on filesystem and started docker
stop on runlevel [!2345]

script
 [ -e /etc/default/weave ] && . /etc/default/weave
 /usr/bin/docker rm -f weave || true
 <%= scope.lookupvar('weave::bin_dir') %>/weave launch <% if @cluster_peers %><%= @cluster_peers %><% end %>
end script

respawn
respawn limit 10 10
kill timeout 10
